# docker/supervisord.conf
; TOON Oyster CVM - supervisord configuration
;
; Manages two processes inside a single Oyster CVM container:
;   1. toon (priority=10): BLS HTTP + Nostr Relay WS + Bootstrap
;   2. attestation (priority=20): TEE attestation HTTP server
;
; The connector is EXTERNAL and not managed by supervisord.
; supervisord runs as root (required to switch users via user= directive),
; but individual programs execute as the non-root 'toon' user.
;
; Log output goes to stdout/stderr with maxbytes=0 (no log rotation)
; because logs are captured externally by the Oyster CVM runtime.

[supervisord]
nodaemon=true
user=root

[program:toon]
priority=10
command=node /app/dist/entrypoint-town.js
user=toon
autorestart=true
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:attestation]
priority=20
command=node /app/dist/attestation-server.js
user=toon
autorestart=true
startsecs=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
